# smf.timer
[Unit]
Description=Fetch SMF recent posts and mail them (timer)

[Timer]
OnBootSec=2min
OnUnitActiveSec=7min

[Install]
WantedBy=timers.target
